# scripts/com.helm.jobs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.helm.jobs</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/user/Projects/helm/feeds/job-applications.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/user/Projects/helm</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/user/Projects/helm/logs/jobs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/Projects/helm/logs/jobs.log</string>
</dict>
</plist>
