# scripts/com.helm.voice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.helm.voice</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/user/Projects/helm/voice-server/.venv/bin/python</string>
        <string>server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/user/Projects/helm/voice-server</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/user/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>KOKORO_DEVICE</key>
        <string>cpu</string>
        <key>WHISPER_DEVICE</key>
        <string>cpu</string>
        <key>HF_HUB_DISABLE_XET</key>
        <string>1</string>
        <key>WAKE_WORD</key>
        <string>off</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/user/Projects/helm/logs/voice.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/Projects/helm/logs/voice.log</string>
</dict>
</plist>
