# scripts/com.helm.tokens.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.helm.tokens</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>__HELM_ROOT__/feeds/claude-tokens.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__HELM_ROOT__</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HELM_ROOT__/logs/tokens.log</string>
    <key>StandardErrorPath</key>
    <string>__HELM_ROOT__/logs/tokens.log</string>
</dict>
</plist>
